# systemd/system/keep-external-disk-alive.service
[Unit]
Description=keep external disk alive
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/touch /mnt/repository/.keepalive
